# board/common/rootfs-overlay/etc/systemd/system/client.service.d/override.conf
[Unit]
OnFailure=satos-emergency.service
Wants=rauc.service
After=rauc.service
