# board/common/rootfs-overlay/etc/systemd/system/data.mount
[Unit]
Description = required data usb stick
DefaultDependencies=no
Conflicts=umount.target
Before=local-fs.target umount.target
OnFailure=satos-emergency.service

[Mount]
What=LABEL=discosat-data
Where=/data
Type=ext4
Options=defaults,noatime
TimeoutSec=15

[Install]
WantedBy=local-fs.target