# board/common/rootfs-overlay/etc/systemd/system/usb-modem.service
[Unit]
Description=USB Modem initializing service

[Service]
Type=simple
ExecStart=/usr/bin/modem_manager

[Install]
WantedBy=multi-user.target